# tomcat7.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/bin/bash -c "/usr/local/tomcat/bin/startup.sh"